# deploy/dpserver
#!/bin/sh
#
# dpserver:   Data Processing Server
#
# chkconfig: 345 98 02
# description:  This is a daemon listens for data processing commands and executes them \
#               with the priviledges of the user requesting the action.
#
# Source function library.
. /etc/rc.d/init.d/functions

# Prepare environment
for i in /cmcf_apps/profile.d/*.sh; do
	if [ -r "$i" ]; then
	    . $i
	fi
done
unset i

# services parameters
servicename='dpserver'
pidfile='/var/run/dpserver.pid'
logfile='/var/log/dpserver.log'
appfile=$DPS_PATH/autoprocess/services/server.py
umask=022

export MPLCONFIGDIR=/tmp

# Sanity checks.
[ -f $appfile ] || exit 0

start() {
    echo -n $"Starting DPM Service: "
    daemon  --pidfile=$pidfile $appfile --log=$logfile &
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch /var/lock/subsys/$servicename
}

stop() {
    echo -n $"Stopping Data Processing Service: "

    killproc -p $pidfile
    RETVAL=$?
    echo
    if [ $RETVAL -eq 0 ]; then
        rm -f /var/lock/subsys/$servicename
        rm -f /var/run/$pidfile
    fi
}

# See how we were called.
case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status -p $pidfile $servicename
        RETVAL=$?
        ;;
    restart)
        stop
	sleep 3
        start
        ;;
    condrestart)
        if [ -f /var/lock/subsys/$servicename ]; then
            stop
	    sleep 3
            start
        fi
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart}"
        ;;
esac
exit $RETVAL
